# contrib/ci-workers/templates/buildbot-worker.service.j2
[Unit]
Description=Buildbot worker
Wants=network.target
After=network.target

[Service]
Type=forking
PIDFile=/home/{{ buildbot_worker_user }}/{{ buildbot_worker_name }}/twistd.pid
WorkingDirectory=/home/{{ buildbot_worker_user }}
<<<<<<< HEAD
ExecStart=/usr/local/bin/buildbot-worker start {{ buildbot_worker_name }}
ExecReload=/usr/local/bin/buildbot-worker restart {{ buildbot_worker_name }}
ExecStop=/usr/local/bin/buildbot-worker stop {{ buildbot_worker_name }}
=======
ExecStart={{ pip_bin_dir }}/buildbot-worker start {{ buildbot_worker_name }}
ExecReload={{ pip_bin_dir }}/buildbot-worker restart {{ buildbot_worker_name }}
ExecStop={{ pip_bin_dir }}/buildbot-worker stop {{ buildbot_worker_name }}
>>>>>>> f0a9a0868276e6a96a50d846ed4ebbd378cda875
Restart=always
User={{ buildbot_worker_user }}

[Install]
WantedBy=multi-user.target
